# root/etc/s6-overlay/s6-rc.d/rahla/run
#!/usr/bin/with-contenv sh

s6-setuidgid abc /app/rahla/bin/karaf run